# notebooklm_queue/deploy/systemd/podcasts-notebooklm-profile-refresh.timer
[Unit]
Description=Periodic NotebookLM queue profile refresh

[Timer]
OnBootSec=3m
OnUnitActiveSec=12m
Persistent=true
RandomizedDelaySec=60
Unit=podcasts-notebooklm-profile-refresh.service

[Install]
WantedBy=timers.target
