# notebooklm_queue/deploy/systemd/podcasts-notebooklm-queue@.service
[Unit]
Description=Drain NotebookLM queue for %i
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/podcasts
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/podcasts/notebooklm-queue/%i.env
ExecStartPre=/usr/bin/install -d -m 0755 -o root -g root /var/lib/podcasts/notebooklm-queue
ExecStart=/opt/podcasts/notebooklm_queue/deploy/bin/notebooklm-queue-drain-show.sh %i
TimeoutStartSec=6h
Nice=10

